# files/etc/supervisor/conf.d/supervisord.conf
[supervisord]
nodaemon = true

[program:bash]
command = /bin/bash
autostart = true
autorestart = true

[program:sshd]
command = /usr/sbin/sshd -D
autostart = true
autorestart = true

[include]
files = \var\supervisor\*.conf